# vagrant/etc/profile.d/env.sh
# set environemnt variables

export MAGE_IS_DEVELOPER_MODE=1
export MAGE_MODE=developer
export COMPOSER_HOME=/server/.cache/composer
export PATH="/usr/local/bin:$PATH"
