# vagrant/etc/profile.d/ps1.sh
# configure a pretty ps1

export PS1='\[\033[0;36m\]\u@\h\[\033[0m\]:\@:\[\033[0;37m\]\w\[\033[0m\]$ '
